# crontab
# Run the backup script daily at midnight
0 0 * * * /app/backup.sh > /proc/1/fd/1 2>/proc/1/fd/2

# Run WhatsApp reminder dispatcher every 10 minutes
*/10 * * * * /app/reminder-cron.sh > /proc/1/fd/1 2>/proc/1/fd/2